# supervisor/async-task.ini
[program:async-task-proxy]
command=/usr/local/bin/gunicorn --paste /etc/lico/async-task.ini --bind :18086 --log-config /etc/lico/lico.logging.d/async-task.ini --workers 1 --threads 4 --worker-class gevent --timeout 3600 --keep-alive 65 --capture-output
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
priority=13
