# supervisor/confluent-proxy.ini
[program:confluent-proxy]
command=/usr/local/bin/gunicorn --paste /etc/lico/confluent-proxy.ini --bind :18081 --log-config /etc/lico/lico.logging.d/confluent-proxy.ini --threads 3 --timeout 3600 --worker-class gevent --keep-alive 65 --log-level info --access-logfile - --error-logfile - --capture-output
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
priority=8
