# supervisor/mail-agent.ini
[program:mail-agent]
command=/usr/local/bin/gunicorn --paste /etc/lico/mail-agent.ini --bind :18091 --log-config /etc/lico/lico.logging.d/mail-agent.ini --threads 3 --worker-class gevent --timeout 3600 --keep-alive 65 --capture-output
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
priority=7
