# supervisor/rabbitmq.ini
[program:rabbitmq]
command=rabbitmq-server
autostart=true
autorestart=true
startsecs=3
startretries=3
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/lico/%(program_name)s.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups=10
priority=4

